# etc/supervisor/conf.d/v2ray.ini
[program:v2ray]
command=/usr/bin/v2ray -config /etc/v2ray/v2ray_config.json
redirect_stderr=true
stdout_logfile=/tmp/v2ray.log
